# scripts/crontab/arch_crontab
# === Crontab pro Arch Linux desktop ===
# V코echny cesty jsou absolutn칤, v칳stup logov치n (pokud je pot콏eba)

# 游대 Denn칤 z치loha pomoc칤 backup.sh ve 20:27
27 20 * * * /home/archos/automatizace/skripty/backup.sh >> /home/archos/cron.log 2>&1

# 游댒 Kontrola aktualizac칤 ka쬯칠 3 hodiny
0 */3 * * * /home/archos/update_notifier.sh

